# systemd/.config/systemd/user/clipmenud.service.d/override.conf
[Service]
Environment=CM_SELECTIONS=clipboard
